# chatd.service
[Unit]
Description=simple chat backend
After=network.target

[Service]
Type=simple
User=chatd
WorkingDirectory=/opt/chatd
EnvironmentFile=/etc/default/chatd
ExecStart=/opt/chatd/venv/bin/gunicorn \
    --bind 0.0.0.0:5001 \
    --worker-class gthread \
    --workers 1 \
    --threads 4 \
    --timeout 0 \
    --keep-alive 75 \
    --graceful-timeout 60 \
    chatd:app
Restart=on-failure
NoNewPrivileges=yes
PrivateTmp=yes

[Install]
WantedBy=multi-user.target
